# root/etc/services.d/jackett/run
#!/usr/bin/with-contenv bash

exec \
  s6-setuidgid abc /app/Jackett/jackett -x -d /config --NoUpdates
